# ficheck.cron
# /etc/cron.d/anacron: crontab entries for the anacron package

SHELL=/bin/sh

30 */2 * * *   root	[ -x /usr/bin/ficheck.py ] && /usr/bin/ficheck.py -ru | /usr/local/bin/mail_stuff.py -s "ALERT: [ficheck] `hostname --fqdn`" root
